# auto_down.cron
*/1 * * * *  /usr/src/OneGuyAPI/auto_down.sh
